# instfiles/openrc/xrdp-sesman
#!/sbin/openrc-run

command="/usr/local/sbin/xrdp-sesman"
pidfile="/run/${RC_SVCNAME}.pid"
