# instfiles/openrc/xrdp-sesman
#!/sbin/openrc-run

command="__BASE__/sbin/xrdp-sesman"
pidfile="/run/${RC_SVCNAME}.pid"
