# instfiles/openrc/xrdp
#!/sbin/openrc-run

command="__BASE__/sbin/xrdp"
pidfile="/run/${RC_SVCNAME}.pid"
