# debian/postrm
#!/bin/sh
update-rc.d -f xrdp_control.sh remove > /dev/null
