# debian/postinst
#!/bin/sh
update-rc.d xrdp_control.sh defaults > /dev/null
