# debian/xrdp.cron.daily
#!/bin/sh

if [ -x /usr/sbin/xrdp-reap-sessions ] ; then
	/usr/sbin/xrdp-reap-sessions
fi
