# instfiles/xrdp-sesman.service.in
[Unit]
Description=xrdp session manager
Documentation=man:xrdp-sesman(8) man:sesman.ini(5)
After=network.target
StopWhenUnneeded=true
BindsTo=xrdp.service

[Service]
Type=exec
EnvironmentFile=-@sysconfdir@/sysconfig/xrdp
EnvironmentFile=-@sysconfdir@/default/xrdp
ExecStart=@sbindir@/xrdp-sesman $SESMAN_OPTIONS --nodaemon
ExecReload=kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
